# deploy/bfrs-api.service
[Unit]
Description=BFRS Events Staffing API (Netlify-functions replacement)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bfrs
Group=bfrs
WorkingDirectory=/opt/bfrs-api/current
EnvironmentFile=/etc/bfrs-api.env
ExecStart=/usr/bin/node /opt/bfrs-api/current/index.js
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

# Hardening — the service only needs to bind to localhost and read its working dir.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/bfrs-api
RestrictAddressFamilies=AF_INET AF_INET6
LockPersonality=true

[Install]
WantedBy=multi-user.target
